# chef/cookbooks/baseconfig/files/default/rc.local
#!/bin/sh
sleep 10
/usr/local/bin/uwsgi --ini /home/ubuntu/project/recipemanager/uwsgi.ini --daemonize /var/log/recipemanager.log
